# .config/systemd/user/btc_ticker.timer
[Unit]
Description=BTC price ticker update timer

[Timer]
OnActiveSec=5min
OnUnitActiveSec=10min
AccuracySec=5min

[Install]
WantedBy=timers.target
