# .config/systemd/user/btc_ticker.service
[Unit]
Description=BTC price ticker
ConditionFileIsExecutable=%h/bin/kraken_dzen

[Service]
Type=oneshot
ExecStart=%h/bin/kraken_dzen %t/btc-ticker
