# .config/systemd/user/de.target
[Unit]
Description=Desktop Environment (graphical-session)
BindsTo=graphical-session.target
After=graphical-session.target
