# .config/systemd/user/graphical.target
# ~/.config/systemd/user/graphical.target
[Unit]
Description=Graphical Desktop
BindsTo=graphical-session.target
Wants=graphical-session-pre.target
