# .config/systemd/user/xmonad.service
# file: ~/.config/systemd/user/xmonad.service
[Unit]
Description=XMonad tiling window manager (xinit)
PartOf=graphical.target

[Service]
ExecStartPre=-%E/xmonad/build
ExecStart=%h/.local/bin/xmonad-my
Restart=always
SyslogIdentifier=%P
StandardError=journal

[Install]
WantedBy=graphical.target
