# .config/systemd/user/dunst.service.d/override.conf

[Install]
WantedBy=
WantedBy=de.target
