# .config/systemd/user/clipmenud.service.d/override.conf
[Service]
Environment=CM_MAX_CLIPS=5000 CM_OWN_CLIPBOARD=1 "CM_SELECTIONS=clipboard" CM_DEBUG=0
#Environment="CM_IGNORE_WINDOW=regex"

Nice=10
RestartSec=60s

[Install]
WantedBy=
WantedBy=de.target
