# .config/systemd/user/mpd-notification.service.d/override.conf
[Service]
NoNewPrivileges=yes
MemoryDenyWriteExecute=yes
RestrictRealtime=yes
Nice=10
