# .config/systemd/user/synergys.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/synergys --no-daemon
